# 8022_nginx_gunicorn_socket_pass_flask/8022_nginx_gunicorn_flask_socket_pass.socket
[Unit]
Description=helloworld socket 12

[Socket]
ListenStream=/run/8022_nginx_gunicorn_flask_socket_pass/socket

[Install]
WantedBy=sockets.target
