# 8001_nginx_simple/helloworld_1_nginx_simple.socket
[Unit]
Description=helloworld socket 1

[Socket]
ListenStream=/run/helloworld_1_nginx_simple/socket

[Install]
WantedBy=sockets.target
