# 8003_gunicorn_socket_global/helloworld_3_gunicorn_socket_global.service
[Unit]
Description=helloworld gunicorn socket daemon
After=network.target
Requires=helloworld_3_gunicorn_socket_global.socket


[Service]
PIDFile=/run/helloworld_3_gunicorn_socket_global/pid
User=www-data
Group=www-data
RuntimeDirectory=helloworld_3_gunicorn_socket_global
WorkingDirectory=/var/helloworld_3_gunicorn_socket_global
ExecStart=gunicorn \
	-w 1 \
	--bind unix:/run/helloworld_3_gunicorn_socket_global/socket \
	--pid /run/helloworld_3_gunicorn_socket_global/pid \
	hello:app
ExecReload=/bin/kill -s HUP $MAINPID
ExecStop=/bin/kill -s TERM $MAINPID
PrivateTmp=true

[Install]
WantedBy=multi-user.target

