# 8002_nginx_socket_pass/8001_helloworld_nginx_socket.socket
[Unit]
Description=helloworld socket 1

[Socket]
ListenStream=/run/8001_helloworld_nginx_socket/socket

[Install]
WantedBy=sockets.target
