# gunicorn_socket_global/helloworld_3_gunicorn_socket_global.socket
[Unit]
Description=helloworld socket 3

[Socket]
ListenStream=/run/helloworld_3_gunicorn_socket_global/socket

[Install]
WantedBy=sockets.target
